# actividad4/file.service
[Unit]
Description=Servicio de saludo y fecha actual

[Service]
Type=simple
ExecStart=/ruta/a/tu/script/script.sh

[Install]
WantedBy=default.target
